# src/whm-root-ip-guard-alert.service
[Unit]
Description=Alert that whm-root-ip-guard failed

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c 'echo "whm-root-ip-guard failed on $(hostname) at $(date)" | mail -s "WHM guard down: $(hostname)" root'
